# utils/logsplitter-monitor-logger.service
[Unit]
Description=LogSplitter Monitor Database Logger
After=network.target

[Service]
Type=simple
User=rayvn
WorkingDirectory=/home/rayvn/monitor/LogSplitter_Monitor/utils
ExecStart=/home/rayvn/monitor/LogSplitter_Monitor/utils/venv/bin/python monitor_database_logger.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monitor-logger

[Install]
WantedBy=multi-user.target
